# tg-cli/kill.cron
*/5 * * * * kill $(pgrep telegram-cli) > /dev/null 2>&1
